# dist_pc/DEBIAN/postinst
#!/bin/sh
invoke-rc.d gfw enable || echo "Success"

