# dist_pc/DEBIAN/postinst
#!/bin/sh
update-rc.d gfw defaults || echo ""
update-rc.d gfw enable || echo ""

