# dist_pc/DEBIAN/postinst
#!/bin/sh
update-rc.d gfw defaults 95 || echo ""
update-rc.d gfw enable || echo "Success"

